# backup/systemd/backup-sugar-sugar.timer
[Unit]
Description=Daily backup for Sugar-Sugar data/input

[Timer]
OnCalendar=*-*-* 03:17:00
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
